# RNG/pi-rng-kiosk/system/pi-rng-kiosk.service
[Unit]
Description=Pi RNG Kiosk
After=graphical-session.target

[Service]
Type=simple
Environment=QT_QPA_PLATFORM=wayland
ExecStart=__PROJECT_DIR__/.venv/bin/python __PROJECT_DIR__/app.py
Restart=on-failure

[Install]
WantedBy=default.target
