# RNG/pi-rng-kiosk/system/pi-rng-kiosk.service
[Unit]
Description=Pi RNG Kiosk
After=graphical-session.target

[Service]
Type=simple
Environment=QT_QPA_PLATFORM=wayland
ExecStart=/bin/bash -lc "cd __PROJECT_DIR__ && ./scripts/update.sh && ./.venv/bin/python app.py"
Restart=on-failure

[Install]
WantedBy=default.target
